# package-name/config/supervisord.conf.template
[supervisord]
nodaemon=true

[program:server]
directory=/code
command=$packageName
autostart=true
autorestart=true
